# service/dep2/finish
#!/bin/sh

echo "[dep2] Finishe dep2"
echo 0 > /tmp/dep2
