# service/dep1/finish
#!/bin/sh

echo "[dep1] Finish dep1"
echo 0 > /tmp/dep1
