# service/test/finish
#!/bin/sh

echo "[test] Finished installation!"
echo
