# service/test/finish
#!/bin/sh

echo "Finished installation!"
echo
